# crontab/crontab

0 0 * * * /usr/local/bin/python3 /geoip/update_database.py
